# dev-db/mysql-community/files/6.0.9_alpha/init/mysqld
#!/sbin/runscript

CONFIGDIR="/etc/mysql"
DATADIR="/var/lib/mysql"
PIDFILE="/var/run/mysqld/mysqld.pid"

depend() {
    need net
    use dns
}

checkconfig() {
    if [ ! -f $CONFIGDIR/my.cnf ] ; then
        eerror "No $CONFIGDIR/my.cnf file exists!"
        return 1
    fi

    if [ ! -d $DATADIR/mysql ] ; then
        eerror "The mysql database is not installed! Please run:"
        eerror "emerge -va dev-db/mysql --config"
        return 1
    fi
}

start() {
    ebegin "Checking mysqld configuration"
        checkconfig || return 1
    eend $?

    ebegin "Starting mysqld"
        mysqld_safe > /dev/null 2>&1 &
    eend $?
}

stop () {
    ebegin "Stopping mysqld"
        start-stop-daemon --stop --quiet --pidfile=$PIDFILE
    eend $?
}
